# cookbooks/bcpc/templates/default/routemon.conf.erb
# Make sure we start before an interface receives traffic
start on (starting networking)

stop on runlevel [!023456]

respawn

console log

exec /usr/local/bin/routemon.pl \
<%= node['bcpc']['routemon']['numfixes']%> \
<%= node['bcpc']['management']['interface']%> \
<%= node['bcpc']['storage']['interface']%> \
<%= node['bcpc']['floating']['interface']%>

