# cookbooks/bcpc/templates/default/upstart-zabbix-server.conf.erb
################################################
#
#              Generated by Chef
#
################################################

description     "Zabbix Server"
author          "S. CANCHON"
start on vip_won
stop on (runlevel [016] or vip_lost)
respawn
expect daemon
exec /usr/local/sbin/zabbix_server -c /usr/local/etc/zabbix_server.conf
